# pensionplanner.service
[Unit]
Description=Retirement Income Planner (Gunicorn)
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/pensionplanner
EnvironmentFile=/opt/pensionplanner/.env
Environment="PATH=/opt/pensionplanner/venv/bin"
ExecStart=/opt/pensionplanner/venv/bin/gunicorn --workers 2 --bind 127.0.0.1:8001 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
